# aggregator-app/src/deb/control/postrm
#!/bin/bash

action=$1
version=$2

if [ "$action" != upgrade ];then
    update-rc.d aggregator remove
fi

exit 0
